# webcam-fix-book5/cam-rot-check-upstream.service
[Unit]
Description=Samsung Galaxy Book — camera rotation ACPI override upstream monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/cam-rot-check-upstream.sh
RemainAfterExit=no
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
